# oc3c-ntp-clock/files/etc/rc.local
/etc/init.d/sysntpd disable

#snmpd -c /sysconf/net-snmp/snmp/snmpd.conf

ntpd -lN

########
# Server
#/SDC/Clock 38400 1 1 1 1 0 60 &

########
# Client
/SDC/Clock 4800 1 1 1 0 1 60 &

exit 0
